# user_monitor.service
[Unit]
Description=usermonitor service
Wants=network.target
After=iotmanager.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/home/sys_run-log_monitor/sys_run-log_monitor.sh
StandardOutput=journal
KillMode=mixed
Restart=always
RestartSec=8s

[Install]
WantedBy=multi-user.target
